# templates/com.vccc.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Example launchd agent generated by scripts/cron_to_launchd.py -->
<!-- This one runs a session every weekday at 09:07 LOCAL TIME (not UTC). -->
<!-- Label convention: com.vccc.<schedule-name> — group all agents by prefix. -->
<plist version="1.0"><dict>
  <key>Label</key><string>com.vccc.daily-standup</string>
  <key>ProgramArguments</key><array>
    <string>/Users/YOU/.claude/ccc-orchestrator/run-scheduled.sh</string>
    <string>daily-standup</string>
  </array>
  <key>StandardOutPath</key><string>/Users/YOU/.claude/ccc-orchestrator/logs/daily-standup.log</string>
  <key>StandardErrorPath</key><string>/Users/YOU/.claude/ccc-orchestrator/logs/daily-standup.err</string>
  <key>StartCalendarInterval</key><array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>7</integer><key>Weekday</key><integer>1</integer></dict>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>7</integer><key>Weekday</key><integer>2</integer></dict>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>7</integer><key>Weekday</key><integer>3</integer></dict>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>7</integer><key>Weekday</key><integer>4</integer></dict>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>7</integer><key>Weekday</key><integer>5</integer></dict>
  </array>
</dict></plist>
